# console/kiosk/systemd/deckport-kiosk.service
[Unit]
Description=Deckport Kiosk Manager
After=network.target
Requires=network.target

[Service]
Type=simple
User=root
ExecStart=/home/jp/deckport.ai/console/scripts/kiosk_manager.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Create log file
ExecStartPre=/usr/bin/touch /var/log/deckport-kiosk.log
ExecStartPre=/usr/bin/chown root:root /var/log/deckport-kiosk.log

[Install]
WantedBy=multi-user.target
