# console/kiosk/systemd/deckport-console.service
[Unit]
Description=Deckport Console Game
After=graphical-session.target network.target
Wants=graphical-session.target

[Service]
Type=simple
User=deckport
Group=deckport
Environment=DISPLAY=:0
WorkingDirectory=/home/jp/deckport.ai/console/build/linux_x86_64
ExecStart=/home/jp/deckport.ai/console/build/linux_x86_64/deckport_console --fullscreen
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
MemoryMax=2G
CPUQuota=80%

[Install]
WantedBy=default.target
