# console/kiosk/deployment_assets/configs/systemd/network-check.service
[Unit]
Description=Deckport Network Connectivity Check
Before=deckport-kiosk.service

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'if ! ping -c 1 8.8.8.8 > /dev/null 2>&1; then systemctl start wifi-portal.service; fi'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
